# systemd/.config/systemd/user/imwheel.service
[Unit]
Description="imwheel"

[Service]
Type=forking
ExecStart=/usr/bin/imwheel

[Install]
WantedBy=graphical-session.target
