# normcap/.config/systemd/user/normcap.service
[Unit]
Description=normcap background mode

[Service]
Type=exec
ExecStart=normcap --background-mode

[Install]
WantedBy=default.target
